# luci-app-natter/root/etc/uci-defaults/luci-natter
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@natter[-1]
	add ucitrack natter
	set ucitrack.@natter[-1].init=natter
	commit ucitrack
EOF

exit 0